# system/user/dotfiles-update.timer
[Unit]
Description=Daily dotfiles auto-update

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
